# etc/systemd/system/boringos-update.service
[Unit]
Description=BoringOS Automated System Update

[Service]
Type=oneshot
ExecStart=/usr/local/bin/boringos-update.py
